# files/lounge_supervisord.conf
[program:lounge]
command=lounge start --home /etc/lounge/
user=lounge
autostart=true
autorestart=true
redirect_stderr=True
