# files/lounge_supervisord.conf
[program:lounge]
command=thelounge start
user=thelounge
autostart=true
autorestart=true
redirect_stderr=True
